# apollo-portal-build/supervisord.conf
[supervisord]
nodaemon=true

[program:apollo-portal]
command=/bin/bash startup.sh
