# apollo-configservice-build/supervisord.conf
[supervisord]
nodaemon=true

[program:apollo-configservice]
command=/bin/bash startup.sh
