# apollo-compose/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/bin/bash start.sh
autorestart=false
