# apollo-adminservice-build/supervisord.conf
[supervisord]
nodaemon=true

[program:apollo-adminservice]
command=/bin/bash startup.sh
